# provisioning/roles/storm/templates/supervisor.conf
[program:storm-{{item}}]
command=storm {{item}}
directory={{install_user_home}}
autorestart=true
user={{install_user}}
